# debian/cron.d
48 17,23,7,12	* * *	news	/usr/lib/gup/gupdate
