# debian/gup.cron.d
48 17,23,7,12	* * *	news	test -x /usr/lib/gup/gupdate && /usr/lib/gup/gupdate
